# mailboot.service.in
[Unit]
Description=Send report mail at boot
After=network.target

[Service]
Type=oneshot
ExecStart=@bindir@/mailboot -s

[Install]
WantedBy=multi-user.target
